# deploy/systemd/fleetd.service
[Unit]
Description=Fleet Terminal backend (fleetd) — browser-based SSH PAM gateway
Documentation=https://github.com/fleet-terminal/fleet-terminal
After=network-online.target
Wants=network-online.target
# Uncomment if Postgres/Redis run on the same host via systemd:
# After=postgresql.service redis.service
# Wants=postgresql.service redis.service

[Service]
Type=simple
# Runtime configuration (DB URL, secrets, public URL, etc.). See README.md.
EnvironmentFile=/etc/fleet/fleet.env
ExecStart=/usr/local/bin/fleetd
Restart=on-failure
RestartSec=5s
TimeoutStopSec=30s

# --- Identity -------------------------------------------------------------
# A persistent, allocation-free system user. fleetd runs unprivileged.
# (Use DynamicUser=yes instead if you do NOT need a stable uid for the
#  recordings directory; with DynamicUser, set StateDirectory below and drop
#  the User/Group lines.)
User=fleet
Group=fleet

# --- Filesystem & state ---------------------------------------------------
# Creates/owns /var/lib/fleet (0750) for session recordings; FLEET_RECORDING_DIR
# should point here (default /var/lib/fleet/recordings).
StateDirectory=fleet
StateDirectoryMode=0750
RuntimeDirectory=fleet
RuntimeDirectoryMode=0750
WorkingDirectory=/var/lib/fleet

# --- Hardening ------------------------------------------------------------
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectHostname=yes
ProtectProc=invisible
ProcSubset=pid
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
RemoveIPC=yes
# Network: TCP/IP only (HTTP API + outbound SSH); no unix-only or packet sockets.
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
# System call allow-list to common, safe groups.
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources
SystemCallArchitectures=native
# Capabilities: none required (binds :8080, a non-privileged port).
CapabilityBoundingSet=
AmbientCapabilities=
UMask=0077

# --- Resource limits ------------------------------------------------------
LimitNOFILE=65536
TasksMax=4096

[Install]
WantedBy=multi-user.target
